# etc/rc.d/tcpbench
#!/bin/ksh

daemon="/usr/bin/tcpbench"
daemon_flags="-s"
daemon_user=_tcpbench

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
